# src/systemd/kiro2chat.service
[Unit]
Description=Kiro2Chat — ACP Chat Bridge
After=network.target

[Service]
Type=simple
WorkingDirectory=__WORKING_DIR__
ExecStart=__UV_PATH__ run kiro2chat run
Restart=on-failure
RestartSec=5
KillMode=control-group
KillSignal=SIGTERM
TimeoutStopSec=10
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
