# kiro2chat@.service
[Unit]
Description=kiro2chat API Gateway
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/repos/kiro2chat
Environment=PATH=/home/%i/.local/bin:/usr/bin
Environment=KIRO_DB_PATH=/home/%i/.local/share/kiro-cli/data.sqlite3
ExecStart=/home/%i/.local/bin/uv run kiro2chat api
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/%i/.config/kiro2chat /home/%i/.local/share/kiro-cli

[Install]
WantedBy=multi-user.target
